# packaging/launchd/com.lookslikecode.ChargeLimitMonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lookslikecode.ChargeLimitMonitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/charge-limit-monitor</string>
    <string>--target</string>
    <string>80</string>
    <string>--verbose</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/charge-limit-monitor.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/charge-limit-monitor.err.log</string>
</dict>
</plist>
